# server/qwopus.service
[Unit]
Description=Qwopus 27B AI Stack
After=network.target docker.service
Wants=docker.service

[Service]
Type=forking
User=ubuntu
ExecStart=/home/ubuntu/QwopusUncensored/server/startup.sh
RemainAfterExit=yes
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
